# files/etc/crontabs/root
#每月定时重启路由器
11 1 1 * * sleep 5 && touch /etc/banner && reboot
#每10分钟检测网络，重启WAN
*/10 * * * * /etc/wan_check.sh
#github-host
13 3 * * * /etc/hostsUpdate.sh
##################################
